# systemd/user/umans-limiter.service
[Unit]
Description=Local UMANS API concurrency limiter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.local/share/umans-local-provider
ExecStart=%h/.local/share/umans-local-provider/bin/bun %h/.local/share/umans-local-provider/src/umans-limiter.ts
Restart=on-failure
RestartSec=2
Environment=UMANS_LIMITER_HOST=127.0.0.1
Environment=UMANS_LIMITER_PORT=8319
Environment=UMANS_LIMITER_MAX_CONCURRENCY=4
Environment=UMANS_LIMITER_UPSTREAM=https://api.code.umans.ai

[Install]
WantedBy=default.target
